# deployment/supervisor/forum-tasks.conf
[program:forum-tasks]
command=/home/ubuntu/forum/venv/bin/rq worker forum-tasks
numprocs=1
directory=/home/ubuntu/forum
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
